# supervisord/supervisord-proxy.conf
[program:danted]
command=/start_danted.sh
numprocs=1
autostart=true
autorestart=true

